# config/RPM/periscoped.service
[Unit]
Description=Persicope Service
After=network.target
Requires=mongod.service redis.service

[Service]
Type=simple
User=periscope
EnvironmentFile=/etc/sysconfig/periscoped
ExecStart=/usr/bin/supervisord -c /etc/supervisor/supervisor.conf

[Install]
WantedBy=multi-user.target